# infra/systemd/hlprivateer-web.service
[Unit]
Description=HL Privateer Web
After=network-online.target
Wants=network-online.target
Requires=hlprivateer-api.service hlprivateer-runtime.service hlprivateer-ws.service
After=hlprivateer-api.service hlprivateer-runtime.service hlprivateer-ws.service

[Service]
Type=simple
User=hlprivateer
Group=hlprivateer
WorkingDirectory=/opt/hlprivateer.xyz
EnvironmentFile=/opt/hlprivateer.xyz/config/.env
LoadCredential=hlprivateer.env:/etc/hlprivateer/credentials/hlprivateer.env
EnvironmentFile=-/run/credentials/%N/hlprivateer.env
ExecStart=/usr/local/bin/bun /opt/hlprivateer.xyz/apps/web/.next/standalone/apps/web/server.js
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateUsers=true
PrivateDevices=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/hlprivateer.xyz /var/log/hlprivateer
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectClock=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=true
MemoryDenyWriteExecute=false
SystemCallFilter=@system-service
SystemCallArchitectures=native
UMask=0077
RestartSteps=3
RestartMaxDelaySec=30
StartLimitBurst=5
StartLimitIntervalSec=60

[Install]
WantedBy=multi-user.target
